# modules/adapter/templates/st2_gunicorn_init/init.service.erb
[Unit]
Description=StackStorm <%= @_subsystem %> gunicorn script.
After=network.target

[Service]
Type=simple
Environment="PYTHONPATH=<%= @_python_pack %>/<%= @_subsystem %>/<%= @_subsystem %>"
ExecStart=gunicorn_pecan <%= @_python_pack %>/<%= @_subsystem %>/gunicorn_config.py \
      -k eventlet -b unix:<%= @socket %> --threads <%= @threads %> \
      --workers <%= @workers %> -u <%= @user %> -g <%= @group %>
TimeoutSec=60
PrivateTmp=true
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
