# modules/adapter/templates/st2_uwsgi_init/init.conf.erb
# StackStorm <%= @_subsystem %> uWSGI script.
#
# Anchors to <%= @_subsystem %> parent init script
# to maintain compatability with st2ctl
#
# This file is managed by Puppet

description     "StackStorm <%= @_subsystem %> uWSGI Daemon"
author          "StackStorm Engineering <opsadmin@stackstorm.com>"

start on runlevel [2345]
stop on runlevel [016]

respawn
respawn limit 2 5

umask 007
kill timeout 60

exec uwsgi --ini /etc/uwsgi.d/<%= @_subsystem %>.ini
